# systemd/codex-telegram-gateway.service
[Unit]
Description=Codex Telegram Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Replace these values with the user/group chosen for your deployment.
User=service-user
Group=service-group

# Replace <project_dir> with the absolute path where the project is installed.
WorkingDirectory=<project_dir>
EnvironmentFile=<project_dir>/.env
ExecStart=<project_dir>/.venv/bin/python3 -m codex_telegram_gateway --config <project_dir>/config.toml --env-file <project_dir>/.env
Restart=always
RestartSec=5
TimeoutStopSec=20
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only

# Replace these paths with the directories your service must read/write:
# - project directory
# - allowed workspace roots
# - runtime/log directories
# - Codex binary/profile locations if they live under protected home directories
ReadWritePaths=<project_dir> <allowed_root_1> <allowed_root_2> <state_dir> <log_dir> <codex_home_dir>

[Install]
WantedBy=multi-user.target
